# rc.local
#!/bin/sh -e

iptables -I INPUT -p tcp --dport 443 -j ACCEPT
/usr/local/bin/xray -c /usr/local/bin/trojan.json &
/usr/bin/caddy run --environ --config /etc/caddy/Caddyfile &

exit 0
